# native-linux/systemd/wave3-daemon.service
[Unit]
Description=Elgato Wave:3 native control daemon
Documentation=See README.md in the source tree
After=pipewire.service wireplumber.service dbus.service
Wants=pipewire.service wireplumber.service

[Service]
Type=simple
ExecStart=%h/.local/bin/wave3-daemon
Restart=on-failure
RestartSec=5

# Environment
Environment="G_MESSAGES_DEBUG=all"

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Required paths
ReadWritePaths=%t
ReadWritePaths=/dev/bus/usb

# USB access for libusb
PrivateDevices=no
DeviceAllow=/dev/bus/usb rw

[Install]
WantedBy=default.target
